# contrib/systemd/osmo-pcap-client.service
[Unit]
Description=PCAP Client for the PCAP aggregation
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
StateDirectory=osmocom
WorkingDirectory=%S/osmocom
ExecStart=/usr/bin/osmo-pcap-client -c /etc/osmocom/osmo-pcap-client.cfg
RestartSec=2

[Install]
WantedBy=multi-user.target

